# home/.config/systemd/user/clash.service
[Unit]
Description=Clash

[Service]
ExecStart=/usr/bin/clash -f ~/.config/clash/config.yaml

[Install]
WantedBy=default.target